# systemd/spring-pet-clinic.service
[Unit]
Description=Spring Pet Clinic

[Service]
ExecStart=/app/spring-pet-clinic/start.sh

[Install]
WantedBy=multi-user.target
